# etc/profile.d/proxy.csh
setenv http_proxy http://138.209.111.74:80
setenv https_proxy http://138.209.111.74:80
setenv ftp_proxy http://138.209.111.74:80
setenv HTTP_PROXY http://138.209.111.74:80
setenv FTP_PROXY http://138.209.111.74:80
