# etc/profile.d/proxy.sh
http_proxy=http://138.209.111.74:80
https_proxy=http://138.209.111.74:80
ftp_proxy=http://138.209.111.74:80
no_proxy=.lmmfc-vrsil.com
HTTP_PROXY=http://138.209.111.74:80
HTTPS_PROXY=http://138.209.111.74:80
FTP_PROXY=http://138.209.111.74:80
export http_proxy https_proxy ftp_proxy no_proxy HTTPS_PROXY HTTP_PROXY FTP_PROXY
